# src/usr/lib/systemd/system/geoipdbupdater.timer
[Unit]
Description=GeoIP database updater timer

[Timer]
OnCalendar=*-*-* 01:47
Persistent=true

[Install]
WantedBy=timers.target
